# rootfs/etc/supervisor.d/go-carbon.ini
[program:go-carbon]
command        = /usr/bin/go-carbon -config /etc/carbon.conf
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
